# rootfs/var/spool/cron/crontabs/root
0 0 * * * /usr/local/hstc/bin/n-backup-conf
